# systemd/tv-off.service
[Unit]
Description=Turn TV off via HDMI-CEC

[Service]
Type=oneshot
ExecStart=/home/pi/Interface-stock/scripts/tv-off.sh
